# backup-frontend.ini
[program:backup-frontend]
directory=/usr/share/nginx/backup-frontend
command=npm run dev
autostart=true
autorestart=true
